# seizures.cron
# m h   dom mon dow     command
*/5 *   *   *   *       ~/venv/bin/python3 ~/manage.py rescue > ~/last_mail.log
*/5 *   *   *   *       ~/venv/bin/python3 ~/manage.py tosqlite > ~/tosqlite.log
